# ansible/roles/technitium/templates/dns-backup.timer.j2
[Unit]
Description=Run Technitium DNS backup nightly

[Timer]
OnCalendar=*-*-* {{ technitium_backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
